# templates/docker.service.j2
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd -H fd://{% for host in docker_hosts %} -H {{ host }}{% endfor %} --containerd=/run/containerd/containerd.sock
